# saltstack/salt/resources/etc/init/dlgtweb.conf
description "A simple web server for testing the delegator web interface"

start on runlevel [2345]
stop on runlevel [!2345]

setuid gator
setgid www-data

exec /bin/bash /etc/gator/startdlgtweb.sh
